# docker/standalone/etc/s6-overlay/s6-rc.d/supervisor/run
#!/command/execlineb -P

# Use with-contenv to ensure environment variables are available
with-contenv

foreground {
    /usr/bin/supervisord -c /etc/supervisord.conf
}
